# config/minidash.service
[Unit]
Description=Minidash
After=network.target

[Service]
Type=simple
Restart=always
ExecStart=/usr/bin/npm run dev
User=brianduff
